# vpn/supervisor/ipsec.conf
[program:ipsec]
command=/usr/sbin/ipsec start --nofork
autorestart=true