# dbhost/supervisor/cron.conf
[program:cron]
command=service cron start
autorestart=true