# dbhost/supervisor/snmpd.conf
[program:snmpd]
command=/etc/init.d/snmpd start
autorestart=true