# dbhost/supervisor/snmpd.conf
[program:mongodb]
command=/etc/init.d/snmpd start
autorestart=true