# dbhost/supervisor/snmpd.conf
[program:snmpd]
command=service snmpd start
autorestart=true