# nagios/supervisor/apache.conf
[program:apache]
command=/usr/sbin/apache2ctl -DFOREGROUND
priority=1
autorestart=true
