# nagios/supervisor/nagios.conf
[program:nagios]
user = nagios
command = /usr/local/nagios/bin/nagios /usr/local/nagios/etc/nagios.cfg
priority=2
autorestart=true
