# gemp18/supervisor/mongodb.conf
[program:mongodb]
command=/usr/bin/mongod --config /etc/mongod.conf
autorestart=true