# mongodb/supervisor/snmpd.conf
[program:snmpd]
command=snmpd -f
autorestart=true